# 5.6/nginx/image-files/etc/supervisor/conf.d/php-fpm.conf
[program:phpfpm]
command=php-fpm --fpm-config /usr/local/etc/php-fpm.conf -c /usr/local/etc/
autostart=true
autorestart=true
autorstart=true
startsecs=5
startretries=100
redirect_stderr=true
