# 7.0/nginx/image-files/etc/supervisor.d/nginx.conf
[program:nginx]
command=nginx
autostart=true
autorestart=true
startsecs=5
startretries=100
redirect_stderr=true